# docker/service_app/services/starman/run
#!/bin/sh
exec chpst -u appuser starman --listen 127.0.0.1:8080 --workers 4 scripts/feedback_app.psgi
