# docker/service_app/services/starman/run
#!/bin/sh
exec chpst -u appuser starman --port 8080 /app/feedback_app.psgi
